# website.conf
[program:website]
command=/home/discord/discord/gunicorn_start
user=discord
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/discord/discord/logs/gunicorn.log

