# website/paypal-discord-bot-ChampNetwork.conf

[program:paypal-discord-bot-ChampNetwork]
command=source /home/discord/venv/bin/activate
command=/home/discord/venv/bin/python3  /home/discord/discord/website/paypal-discord-bot-ChampNetwork.py
user=discord
stdout_logfile=/home/discord/discord/logs/CHAMP_bot_process.log
redirect_stderr=true
